# SERVER_CONFIG/ubuntu_init/unicorn
#! /bin/bash
### BEGIN INIT INFO
# Provides:          zippelback
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: unicorn handler for zippelback
# Description:       unicorn handler for zippelback
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

pidfile='/var/run/nazgul.pid'
config='/var/www/nazgul/config/unicorn.rb'
name='nazgul'

# . "/usr/local/rvm/scripts/rvm"
. /usr/local/lib/rvm

case "$1" in
  start)
	echo -n "Starting $name: "
	unicorn -c $config -E production -D
	echo "done."
	;;
  stop)
	echo -n "Stopping $name: "
	cat $pidfile | xargs kill -QUIT
    rm -f ${pidfile}
	echo "done."
	;;
  restart)
	echo -n "Restarting $name: "
	cat $pidfile | xargs kill -QUIT
    rm -f ${pidfile}
	sleep 1
	unicorn -c $config -E production -D
	echo "done."
	;;
  reload)
	echo -n "Reloading $name: "
	cat $pidfile | xargs kill -USR2
	echo "done."
	;;
  status)
	if [ -f $pidfile ]
	then
		echo "* $name is running"
	else
		echo "* $name is not running"
	fi
	;;
  *)
	echo "Usage: service nazgul {start|stop|restart|reload|status}" >&2
	exit 1
	;;
esac

exit 0